# deeplearning4j-scaleout/deeplearning4j-scaleout-akka/bin/supervisord-deeplearning-master.conf
[supervisord]
childlogdir=/var/log/supervisord/
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[inet_http_server] 
port=0.0.0.0:9001

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
prompt = >>

[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
autorestart=true
stopsignal=KILL
stdout_logfile=/var/log/sendalyzeit/zookeeper.log
stderr_logfile=/var/log/sendalyzeit/zookeeper-error.log


[program:training]
command=/root/java-deeplearning/deeplearning-parent/bin/